# src/php-4.4/build/modules/build/files/etc/supervisor/conf.d/php.conf
[program:php]
command = /usr/local/src/phpfarm/inst/php-4.4.9/bin/php -b 0.0.0.0:8002
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
